# rootfs/etc/systemd/system/gdm.service.d/override.conf
[Service]
ExecStopPost=-/bin/sh -c 'cp /var/home/$$(id -nu 1000)/.config/monitors.xml /var/lib/gdm/.config/monitors.xml'
